# docs/launchd/ai.k2lab.gitea-watcher.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.k2lab.gitea-watcher</string>
    <key>ProgramArguments</key>
    <array>
      <string>__REPO_ROOT__/bin/gitea-watcher</string>
      <string>--once</string>
    </array>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>EnvironmentVariables</key>
    <dict>
      <key>COURT_ROOT</key>
      <string>__COURT_ROOT__</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/gitea-watcher.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/gitea-watcher.err.log</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
